# deploy/as400-data-importer-pro.service
[Unit]
Description=AS400 Data Importer Pro
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/as400-data-importer-pro/current
EnvironmentFile=/etc/as400-data-importer-pro.env
ExecStart=/opt/as400-data-importer-pro/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
User=as400pro
Group=as400pro

[Install]
WantedBy=multi-user.target
